# resources/systemd/busboy-collector.service
[Unit]
Description=Fills Busboy database

[Service]
Type=simple
ExecStart=/usr/bin/env python3 /home/noel/busboy/busboy.py
