# resources/systemd/busboy-collector.service
[Unit]
Description=Fills Busboy database

[Service]
Type=simple
ExecStart=python3 /home/noel/busboy/busboy.py
