# resources/systemd/busboy-collector.timer
[Unit]
Description=Starts Busboy database collector in the morning

[Timer]
OnCalendar=06:00
