# resources/systemd/busboy-collector-stop.timer
[Unit]
Description=Stops Busboy database collector at night

[Timer]
OnCalendar=01:00
