# resources/systemd/busboy-collector-stop.service
[Unit]
Description=Stops filling Busboy database

[Service]
Type=oneshot
ExecStart=systemctl stop busboy-collector.service
